# scripts/stats/com.dealerinspire.sbm-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.dealerinspire.sbm-refresh</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/nathanhart/auto-sbm/scripts/stats/run_refresh_metadata.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>6</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/nathanhart/auto-sbm/logs/refresh_metadata.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/nathanhart/auto-sbm/logs/refresh_metadata.err</string>
    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
